# files/etc/uci-defaults/uci-settings
uci set system.@system[0].timezone=CST-8
uci set system.@system[0].zonename=Asia/Shanghai
uci set system.@system[0].description='Xiaomi AX6 With OpenWrt'
uci commit system

uci commit
